# sample.cron
#!/bin/bash
ruby /etc/mysql_backup/mysql_backup.rb backup > /dev/null
exit 0
